# systemd/codex-telegram-bot.service.template
[Unit]
Description=Codex Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKDIR}}
ExecStart={{ENTRYPOINT}} --config-dir {{CONFIG_DIR}}
Restart=always
RestartSec=3

[Install]
WantedBy={{WANTED_BY}}
