# systemd/codex-telegram-bot-daemon.service.template
[Unit]
Description=Codex Telegram Bot - Background Agent Daemon
After=network.target
Wants=codex-telegram-bot.service

[Service]
Type=simple
WorkingDirectory={{WORKDIR}}
ExecStart={{ENTRYPOINT}} --daemon --config-dir {{CONFIG_DIR}}
Restart=always
RestartSec=5
WatchdogSec=120

Environment=CRON_AGENT_HEALTH_FILE=%t/codex-telegram-bot-daemon.health
Environment=CRON_AGENT_HEARTBEAT_SEC=1800
Environment=CRON_AGENT_SYSTEM_WATCH_SEC=300
Environment=CRON_AGENT_GIT_WATCH=1
Environment=CRON_AGENT_DISK_WATCH=1

[Install]
WantedBy={{WANTED_BY}}
